# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:synapse]
command=/start.py

[program:proxy]
command=/proxy/proxy
