# v1.64/systemd-with-workers/system/matrix-synapse.target
[Unit]
Description=Synapse parent target
After=network.target

[Install]
WantedBy=multi-user.target
