# docker/complement/conf/postgres.supervisord.conf
[program:postgres]
command=/usr/local/bin/prefix-log gosu postgres postgres

# Only start if START_POSTGRES=1
autostart=%(ENV_START_POSTGRES)s

# Lower priority number = starts first
priority=1

autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Use 'Fast Shutdown' mode which aborts current transactions and closes connections quickly.
# (Default (TERM) is 'Smart Shutdown' which stops accepting new connections but
#  lets existing connections close gracefully.)
stopsignal=INT
