# docs/systemd-with-workers/system/matrix-synapse.target
[Unit]
Description=Synapse parent target
After=network.target
AllowIsolate=no

[Install]
WantedBy=multi-user.target
